# assignment4/tarecho@.service
[Unit]
Description=Tarecho service

[Service]
ExecStart=~/hw4/compile_troff.py
StandardInput=socket
