# assignment4/tarecho.socket
[Unit]
Description=sad

[Socket]
ListenStream=/run/tarecho
Accept=yes

[Install]
WantedBy=sockets.target
